# dockerfiles/vis-services/supervisord.conf
[supervisord]
nodaemon=true
logfile_maxbytes=5MB
logfile_backups=10
loglevel=debug

[program:gunicorn]
command=gunicorn -c gunicorn.conf.py wsgi:application
directory=/vis-services
startsecs=0
