# dockerfiles/bumblebee/supervisord.conf
[supervisord]
nodaemon=true

[program:bumblebee]
command=python app.py
directory=/vagrant

[program:nginx]
command=/usr/sbin/nginx
startsecs=0
