# dockerfiles/bumblebee/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx
startsecs=0
