# dockerfiles/bumblebee/supervisord.conf
[supervisord]
nodaemon=true

[program:bumblebee]
command=grunt server
directory=/bumblebee