# dockerfiles/api_with_bumblebee/supervisord.conf
[supervisord]
nodaemon=true

[program:gunicorn]
command=python wsgi.py
directory=/adsws

[program:nginx]
command=/usr/sbin/nginx -c /etc/nginx/nginx.conf
startsecs=0
