# dockerfiles/zookeeper/supervisord.conf
[supervisord]
nodaemon=true

[program:zookeeper]
command=/zk.sh
user=vagrant
